# deploy/systemd/leon-daily-learning-sync.timer
[Unit]
Description=Daily LEON Learning Synchronization Timer
Documentation=leon-daily-learning skill

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
